# deploy/travel-log-postgres.container
# deploy/postgres.container
[Unit]
Description=Travel Log PostgreSQL Database
After=travel-log-pod.service travel-log-volume.service
Requires=travel-log-pod.service travel-log-volume.service 

[Container]
Image=docker.io/library/postgres:16-alpine
EnvironmentFile=/etc/containers/app.env
Volume=travel-log-db-data:/var/lib/postgresql/data:Z
Pod=travel-log.pod

# This links our shared volume to the internal Postgres socket directory
Volume=travel-log-postgres-socket:/var/run/postgresql

